# 2021/CameraNormalis/systemd/cn-watcher.path
[Path]
PathModified=/home/pi/dvorky/2021/CameraNormalis/systemd

[Install]
WantedBy=multi-user.target